# deploy/systemd/keiba-api.service
[Unit]
Description=Keiba Yosou FastAPI Server
After=network.target postgresql.service

[Service]
Type=simple
User=ec2-user
Group=ec2-user
WorkingDirectory=/home/ec2-user/keiba-yosou
Environment="PATH=/home/ec2-user/keiba-yosou/venv/bin"
EnvironmentFile=/home/ec2-user/keiba-yosou/.env
ExecStart=/home/ec2-user/keiba-yosou/venv/bin/uvicorn src.api.main:app --host 0.0.0.0 --port 8000
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
